# packaging/systemd/sonarium-dlna.service
[Unit]
Description=Sonarium DLNA/UPnP MediaServer
Documentation=https://github.com/kidoz/sonarium
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=sonarium
Group=sonarium
EnvironmentFile=/etc/sonarium/sonarium.env
ExecStart=/usr/local/bin/sonarium-dlna
Restart=on-failure
RestartSec=5

# Hardening. Media is read-only; the service writes nothing.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadOnlyPaths=/srv/music

[Install]
WantedBy=multi-user.target
