# packaging/systemd/sonarium-server.service
[Unit]
Description=Sonarium HLS audio server
Documentation=https://github.com/kidoz/sonarium
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=exec
User=sonarium
Group=sonarium
EnvironmentFile=/etc/sonarium/sonarium.env
ExecStart=/usr/local/bin/sonarium-server
Restart=on-failure
RestartSec=5
# Exit 2 is refuse-to-start on bad configuration — restarting cannot help.
# Exit 3 (catalog unreachable at boot) keeps retrying while Postgres starts.
RestartPreventExitStatus=2

# The segment cache is the only writable location (matches
# SONARIUM_HLS_CACHE_DIR=/var/cache/sonarium/hls in the env file).
CacheDirectory=sonarium/hls

# Hardening. ffmpeg segmentation children run under the same sandbox;
# @system-service permits the fork/exec they need. NOTE: a library under
# /home needs ProtectHome=read-only (or a BindReadOnlyPaths= entry) instead
# of ProtectHome=true.
NoNewPrivileges=true
CapabilityBoundingSet=
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
LockPersonality=true
RemoveIPC=true
SystemCallFilter=@system-service
UMask=0077

[Install]
WantedBy=multi-user.target
