# packaging/systemd/sonarium-dlna.service
[Unit]
Description=Sonarium DLNA/UPnP MediaServer
Documentation=https://github.com/kidoz/sonarium
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=exec
User=sonarium
Group=sonarium
EnvironmentFile=/etc/sonarium/sonarium.env
ExecStart=/usr/local/bin/sonarium-dlna
Restart=on-failure
RestartSec=5
# Exit 2 is refuse-to-start on bad configuration — restarting cannot help.
# Exit 3 (catalog unreachable at boot) keeps retrying while Postgres starts.
RestartPreventExitStatus=2

# Hardening. The service reads the media library and talks to Postgres; it
# writes nothing. ProtectSystem=strict already makes the whole filesystem
# read-only. NOTE: a library under /home needs ProtectHome=read-only (or a
# BindReadOnlyPaths= entry) instead of ProtectHome=true.
NoNewPrivileges=true
CapabilityBoundingSet=
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
# AF_NETLINK: LAN-interface auto-detection for the SSDP LOCATION host.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
LockPersonality=true
RemoveIPC=true
SystemCallFilter=@system-service
UMask=0077

[Install]
WantedBy=multi-user.target
